# deployment/systemd/powermap.service
[Unit]
Description=PowerMap unified private-network access
Documentation=https://github.com/steven-ld/PowerMap/tree/main/deployment
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Domain mappings update /etc/hosts and bind loopback HTTPS port 443. The service therefore
# runs as root. Keep the management listener restricted to trusted local access,
# including when web_bind is loopback.
User=root
Group=root
WorkingDirectory=/var/lib/powermap
Environment=RUST_LOG=info
ExecStart=/usr/local/bin/powermap --config /etc/powermap/powermap.toml
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
# Hosts updates use an atomic replace and a cooperative sibling lock, which require writes in /etc.
ReadWritePaths=/etc /var/lib/powermap /var/log/powermap

[Install]
WantedBy=multi-user.target
